# systemd/greetd.service.d/after-nebula-drm.conf
[Unit]
After=nebula-drm.service
Requires=nebula-drm.service
