# systemd/nebula-drm.service
[Unit]
Description=DRM color temperature (nebula-drm)
Documentation=https://github.com/rmarquis/nebula-drm
# Run before the display manager so DRM master is available
Before=greetd.service
# Ensure DRM devices are settled before we try to open them
After=systemd-udev-settle.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/nebula-drm.conf
ExecStart=/usr/bin/nebula-drm -l ${NEBULA_LOCATION} -t ${NEBULA_TEMP} -b ${NEBULA_BRIGHTNESS} ${NEBULA_CARD} -v
# Requires root to acquire DRM master outside compositor
User=root
# Sandbox: only needs /dev/dri access
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
DeviceAllow=/dev/dri rw
DevicePolicy=closed

[Install]
WantedBy=graphical.target
